# ansible/roles/backups/files/backup-docs.service
[Unit]
Description=Backup Docs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-docs.sh
Environment=SSH_AUTH_SOCK=%t/ssh-agent.socket
