# ansible/roles/backups/files/backup-notes.service
[Unit]
Description=Backup WebDAV Notes

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-notes.sh
