# ansible/roles/backups/files/backup-notes.timer
[Unit]
Description=Backup WebDAV notes daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
