# systemd/user/backup-docs.service
[Unit]
Description=Backup Docs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-docs.sh
